# build/projects/AX630C_emmc_arm64_k419_sipeed_nanokvm/rootfs/etc/systemd/system/sysdev.service
[Unit]
Description=System Device Services

[Service]
Type=simple
ExecStart=/opt/scripts/sysdev.sh
ExecStop=/bin/kill -s TERM $MAINPID 
KillSignal=SIGTERM
TimeoutStopSec=3
User=root
Group=root
Restart=always
RestartSec=1
StandardOutput=journal
StandardError=journal
KillMode=mixed

[Install]
WantedBy=multi-user.target
