# build/projects/AX630C_emmc_arm64_k419_sipeed_nanokvm/rootfs/etc/systemd/system/usb-gadget.service
[Unit]
Description=Run usb gadget
DefaultDependencies=no
After=boot.mount
Requires=boot.mount

[Service]
Type=forking
ExecStart=/opt/scripts/usb-gadget.sh start
ExecStop=/opt/scripts/usb-gadget.sh stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
