# build/projects/AX630C_emmc_arm64_k419_sipeed_maixcam2/rootfs/etc/systemd/system/udhcpc.service
[Unit]
Description=Run udhcpc service
DefaultDependencies=no
After=boot.mount
Requires=boot.mount wifi.service

[Service]
Type=simple
ExecStart=/sbin/udhcpc -i wlan0 -s /usr/share/udhcpc/default.script -f -p /run/udhcpc.wlan0.pid
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
