# build/projects/AX630C_emmc_arm64_k419_sipeed_maixcam2/rootfs/etc/systemd/system/usb-gadget.service
[Unit]
Description=Run usb gadget
DefaultDependencies=no
After=boot.mount
Requires=boot.mount

[Service]
Type=forking
PIDFile=/var/run/udhcpd.usb0.pid
ExecStart=/opt/scripts/usb-gadget.sh start
ExecStop=/opt/scripts/usb-gadget.sh stop
Restart=always

[Install]
WantedBy=multi-user.target
