# build/projects/AX630C_emmc_arm64_k419_sipeed_maixcam2/rootfs/etc/systemd/system/usb.service
[Unit]
Description=Run usb
DefaultDependencies=no
After=boot.mount
Requires=boot.mount

[Service]
Type=simple
ExecStart=/opt/scripts/usb.sh start
ExecStop=/opt/scripts/usb.sh stop
Restart=no
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
